# system/etc/profile.d/aliases.sh
# ALIASES:
alias lst='busybox ls -1ac --color=tty'
alias ls='busybox ls -1a --color=tty --group-directories-first'

alias ..='cd ..'
alias ...='.. && ..'
alias ....='... && ..'
alias .....='.... && ..'
alias ......='..... && ..'

alias cp='cp -v -r'
alias mv='mv -v -f'
alias rm='yes | rm -r -i'
