# system/etc/profile.d/paths.sh

# PATHS:
PATH+=':/data/data/com.termux/files/usr/bin'
PATH+=':/data/data/com.termux/files/usr/sbin'
PATH+=':/data/data/com.termux/files/home/.local/bin'
PATH+=':/data/data/com.termux/files/home/go/bin'
